# simple.conf
[program:myserver]
command=gunicorn run:app -w 4 -b 0.0.0.0:80 --timeout 300


[supervisord]
logfile=/root/supervisord.log
loglevel=debug
user=root
